# setup/scripts/ob-celery.conf
start on runlevel [2345]
stop on runlevel [016]
respawn
setuid max
env PATH=/usr/bin/
chdir /var/www/OpenBikes/
exec celery -A tasks worker -B
