# setup/scripts/ob-restart.conf
start on runlevel [2345]
stop on runlevel [016]
respawn
setuid max
*/3600 * * * * service ob-collect restart
